# launchd/dev.phobosclaw.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.phobosclaw</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/robert/.nvm/versions/node/v24.14.0/bin/node</string>
        <string>/Users/robert/phobosclaw/dist/index.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/robert/phobosclaw</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/robert/.nvm/versions/node/v24.14.0/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/robert</string>
        <key>ASSISTANT_NAME</key>
        <string>Rove</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/robert/phobosclaw/logs/phobosclaw.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/robert/phobosclaw/logs/phobosclaw.error.log</string>
</dict>
</plist>
